# Three_task/2/neshcheret_pkg/DEBIAN/postinst
#!/bin/bash

set -e

systemctl daemon-reload

systemctl enable 10_task.service
systemctl start 10_task.service
